# scripts/debian/bitcore-bitcoind.conf
description "Bitcoin Core for OzbCore"
author "BitPay, Inc."

limit nofile 20000 30000

start on runlevel [2345]
stop on runlevel [016]

kill timeout 300
kill signal SIGINT

# user/group for ozbcore daemon to run as
setuid ozbcore
setgid ozbcore

# home dir of the ozbcore daemon user
env HOME=/home/ozbcore

respawn
respawn limit 5 15

script
  exec ozbcore-ozbcored -datadir=/home/ozbcore/.ozbcore/data/
end script
